# ops/cronostudio-ansible-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.cronostudio.ansible-tunnel</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>/Users/nakto/.local/bin/cronostudio-tunnel-daemon.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/nakto/Library/Logs/cronostudio-ansible-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nakto/Library/Logs/cronostudio-ansible-tunnel.err.log</string>
  </dict>
</plist>
